# docs/zsazsa.service.template
# Copy this file to /etc/systemd/system/zsazsa.service and adapt the paths.
# Example:
#   sudo cp docs/zsazsa.service.template /etc/systemd/system/zsazsa.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now zsazsa.service
#
# Production note: when running behind an Apache reverse proxy, edit
# run_webapp.py to bind to 127.0.0.1 instead of 0.0.0.0 so the app
# is not directly reachable on port 5000 from the network.
# See the "Production deployment behind Apache" section in README.md.

[Unit]
Description=zsazsa CTI web application
After=network.target

[Service]
Type=simple
User=CHANGE_ME_USER
Group=CHANGE_ME_GROUP
WorkingDirectory=/opt/zsazsa
Environment=PYTHONUNBUFFERED=1
# SECRET_KEY can also be supplied here instead of in config/__init__.py:
# Environment=SECRET_KEY=change-me-to-a-long-random-string
ExecStart=/opt/zsazsa/venv/bin/python /opt/zsazsa/run_webapp.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
